# dist/systemd/replicator.service
[Unit]
Description=Replicator
Documentation=https://github.com/elsevier-core-engineering/replicator

[Service]
ExecStart=/usr/local/bin/replicator agent -config=/etc/replicator.d/
SuccessExitStatus=13
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
